# config/supervisor.d/apache.conf
[program:apache]
command=apache2ctl -D FOREGROUND -k start
stdout_events_enabled = true
stderr_events_enabled = true
